# sys-apps/exynos5-hwcomposer/files/exynos5-hwcomposer
#!/sbin/runscript
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

PIDFILE=/var/run/${SVCNAME}.pid
BINARY=/usr/bin/${SVCNAME}

depend() {
        use logger
}

start() {
	ebegin "Starting ${SVCNAME}"
	start-stop-daemon --start --background --exec "${BINARY}" -p "${PIDFILE}"
	eend $?
}

stop() {
	ebegin "Stopping ${SVCNAME}"
	start-stop-daemon --stop --exec "${BINARY}" -p "${PIDFILE}" --quiet
	eend $?
}

